# graphite_api/files/graphite-api.conf
description "Graphite-API server"
start on runlevel [2345]
stop on runlevel [!2345]
console log
respawn
exec gunicorn -b 127.0.0.1:8888 -w 4 graphite_api.app:app
